# hello-world.service
[Service]
WorkingDirectory=/home/pi/hello-world
ExecStart=/usr/local/bin/npm start
Restart=always
StandardOutput=syslog
SyslogIdentifier=hello-world
Environment=NODE_ENV=production
# Environment=ANOTHER_KEY=true

[Install]
WantedBy=multi-user.target
